# rc.local
#!/bin/bash
# rc.local - Executed at boot (after network startup)

# Or execute commands directly
python /home/pi/DS18B20_YF-S201_RaspberryPi/DS18B20/main.py &
python /home/pi/DS18B20_YF-S201_RaspberryPi/YF-S201/main.py &
exit 0